# configuration/init/beaglebone/bbb-proxy.conf
# bbb_proxy

start on stopped rc RUNLEVEL=[2345]
start on started roscore

exec /usr/local/bin/start-bbb-proxy
