# configuration/init/beaglebone/beagleboard-udp.service
[Unit]
Description=beagleboard-udp script launcher

[Service]
ExecStart=/usr/local/bin/start-beagleboard-udp

[Install]
WantedBy=multi-user.target
