# configuration/init/beaglebone/beagleboard-udp.conf
# msl_beagle_board_black_udp

start on stopped rc RUNLEVEL=[2345]

exec /usr/local/bin/start-beagleboard-udp
