# configuration/init/beaglebone/beagleboard.conf
# bbb_proxy

start on stopped rc RUNLEVEL=[2345]
start on started bbb-proxy

exec /usr/local/bin/start-beagleboard
