# configuration/init/msl_udp_proxy.conf
# msl_udp_proxy

start on starting roscore
stop on stopped roscore

exec /usr/local/bin/start-udp-proxy



